# auto_management/etc/cron.d/upgrade
0 0 * * * root /bin/bash /etc/upgrade.sh
