# auto_management/etc/cron.d/upgrade
0 0 * * * root /bin/bash /usr/bin/upgrade.sh
